# linux/.config/systemd/user/plasma-kwin_x11.service.d/shorter_TimeoutStopSec.conf
# prevent kwin from hanging for too long during system shutdown
[Service]
TimeoutStopSec=10s
